# crontab
*/20 * * * * /usr/local/bin/python /app/cfddns.py >> /var/log/cron.log 2>&1
